# build/supervisord_dev.conf
[include]
files = /etc/supervisord_common.conf

# esto es para que en entorno de desarrollo haga restart de celery cuando detecta cambios de codigo
[program:watchmedo]
command=watchmedo shell-command --pattern=*.py --recursive --command='/usr/bin/supervisorctl restart celery' /app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock