# build/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0700
chown= nobody:nogroup

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/supervisord.pid
identifier = supervisor


# http://veithen.github.io/2015/01/08/supervisord-redirecting-stdout.html
[program:celery]
command=/wait-for-it.sh --host=app --port=9000 --timeout=500 -- celery -A videocursos worker --loglevel=info
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface